# fireinitiateafterkill.service
[Unit]
Description=chrome

[Service]
Type=simple
Environment=DISPLAY=:0

ExecStart=/home/a/vmware/firefox/firefox -P Default
